# brace/usr/lib/systemd/system/rngd.service.d/99-brace.conf
[Service]
# Hardening
#CapabilityBoundingSet=CAP_SYS_ADMIN
LockPersonality=true
MemoryDenyWriteExecute=true
ProtectProc=invisible
PrivateTmp=yes
#ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
#ProtectKernelModules=true
#ProtectKernelTunables=true
ProtectSystem=strict
RemoveIPC=true
RestrictAddressFamilies=~AF_INET
RestrictAddressFamilies=~AF_INET6
RestrictAddressFamilies=~AF_NETLINK
RestrictAddressFamilies=~AF_PACKET
RestrictAddressFamilies=~AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
UMask=0077
