# brace/usr/lib/systemd/system/ModemManager.service.d/99-brace.conf
[Service]
# Hardening
CapabilityBoundingSet=CAP_SYS_ADMIN
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
#ProtectProc=invisible
PrivateTmp=yes
PrivateUsers=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
#ProtectKernelTunables=true
ProtectSystem=strict
RemoveIPC=true
RestrictAddressFamilies=AF_NETLINK AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
UMask=0077
