# brace/usr/lib/systemd/system/murmur.service.d/90-quickstart.conf
[Service]
ExecStartPre=/bin/sleep 60
#Service can start too quickly before all network interfaces are fully up
