# brace/usr/lib/systemd/system/murmur.service.d/99-brace.conf
[Service]
# Hardening
CapabilityBoundingSet=""
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=yes
PrivateUsers=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadOnlyPaths=/etc/murmur
ReadOnlyPaths=/etc/mumble-server.ini
ReadWritePaths=/var/lib/mumble-server
RemoveIPC=true
RestrictAddressFamilies=~AF_PACKET
RestrictAddressFamilies=~AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged
UMask=0077
