# brace/usr/lib/systemd/system/chronyd.service.d/99-brace.conf
[Service]
#Unset first to allow setting again
ExecStart=
#Use our config instead
ExecStart=/usr/sbin/chronyd -n $OPTIONS -f /etc/chrony.brace.conf
