# brace/usr/lib/systemd/system/wpa_supplicant.service.d/99-brace.conf
[Service]
# Hardening
LockPersonality=true
MemoryDenyWriteExecute=true
#ProtectProc=invisible
PrivateTmp=yes
ProtectControlGroups=true
ProtectHome=read-only
ProtectHostname=true
ProtectKernelLogs=true
#ProtectKernelModules=true
ProtectSystem=strict
ReadOnlyPaths=-/etc/sysconfig/wpa_supplicant
ReadOnlyPaths=-/etc/network
ReadOnlyPaths=-/etc/ifplugd
ReadOnlyPaths=/etc/wpa_supplicant
ReadOnlyPaths=-/home
RemoveIPC=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged
SystemCallFilter=~@resources
UMask=0077
