# brace/usr/lib/systemd/system/libvirtd.service.d/99-brace.conf
[Service]
#workaround hardened_malloc incompatibility with virtiofs
InaccessiblePaths=-/etc/ld.so.preload
