# brace/usr/lib/systemd/system/jellyfin.service.d/99-brace.conf
[Service]
# Hardening
CapabilityBoundingSet=""
LockPersonality=true
NoNewPrivileges=true
PrivateDevices=true
ProtectProc=invisible
PrivateTmp=yes
PrivateUsers=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadOnlyPaths=/etc/sysconfig/jellyfin
ReadOnlyPaths=/usr/share/jellyfin-web
ReadWritePaths=/etc/jellyfin
ReadWritePaths=/var/lib/jellyfin
ReadWritePaths=/var/log/jellyfin
ReadWritePaths=/var/cache/jellyfin
RemoveIPC=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
