# brace/usr/lib/systemd/system/irqbalance.service.d/99-brace.conf
[Service]
# Hardening
#CapabilityBoundingSet="CAP_SETPCAP"
LockPersonality=true
MemoryDenyWriteExecute=true
#NoNewPrivileges=true
PrivateDevices=true
#ProtectProc=invisible
PrivateTmp=yes
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadOnlyPaths=-/etc/default/irqbalance
ReadOnlyPaths=-/etc/sysconfig/irqbalance
ReadOnlyPaths=-/etc/irqbalance
ReadWritePaths=/proc/irq
ReadWritePaths=-/run/irqbalance
ReadWritePaths=-/var/run/irqbalance
RemoveIPC=true
RestrictAddressFamilies=~AF_INET
RestrictAddressFamilies=~AF_INET6
#RestrictAddressFamilies=~AF_NETLINK
RestrictAddressFamilies=~AF_PACKET
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
UMask=0077
