# brace/usr/lib/systemd/system/radiusd.service.d/99-brace.conf
[Service]
# Hardening
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=yes
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadOnlyPaths=/etc/raddb
ReadWritePaths=/var/lib/radiusd
ReadWritePaths=/var/log/radius
ReadWritePaths=/var/run/radiusd
ReadWritePaths=/run/radiusd
RemoveIPC=true
RestrictAddressFamilies=~AF_NETLINK
RestrictAddressFamilies=~AF_PACKET
RestrictAddressFamilies=~AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0077
