# brace/usr/lib/systemd/user/restic-backup@.timer
[Unit]
Description=backup with restic daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=60m
Persistent=true

[Install]
WantedBy=timers.target
