# brace/usr/lib/systemd/system/minetest@.service.d/99-brace.conf
[Service]
# Hardening
CapabilityBoundingSet=""
LockPersonality=true
NoNewPrivileges=true
PrivateDevices=true
#ProtectProc=invisible
PrivateTmp=yes
PrivateUsers=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadOnlyPaths=/etc/minetest
ReadOnlyPaths=/etc/sysconfig/minetest
ReadWritePaths=/var/lib/minetest
RemoveIPC=true
RestrictAddressFamilies=~AF_NETLINK
RestrictAddressFamilies=~AF_PACKET
RestrictAddressFamilies=~AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged
SystemCallFilter=~@resources
UMask=0077
