# brace/usr/lib/systemd/system/php-fpm.service.d/99-brace.conf
[Service]
# Hardening
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=true
#ProtectProc=invisible
PrivateTmp=yes
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadOnlyPaths=-/etc/php-fpm.conf
ReadOnlyPaths=-/etc/php.ini
ReadOnlyPaths=-/etc/php.d
ReadOnlyPaths=-/etc/php-fpm.d
ReadOnlyPaths=-/etc/php-zts.d
ReadWritePaths=/var/www
ReadWritePaths=/var/lib/php
ReadWritePaths=/var/log/php-fpm
ReadWritePaths=-/var/run/php-fpm
ReadWritePaths=-/run/php-fpm
RemoveIPC=true
RestrictAddressFamilies=~AF_PACKET
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@resources
UMask=0077

#workaround hardened_malloc incompatibility
InaccessiblePaths=-/etc/ld.so.preload
